# orderbook-service.service
[Unit]
Description=Hyperliquid Orderbook Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/node
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/bin/bash -c 'docker ps | grep hyperliquid-node-1 || (echo "Hyperliquid node not running" && exit 1)'
ExecStart=/usr/bin/python3 /home/ubuntu/node/integrated_orderbook_service.py --config /home/ubuntu/node/orderbook_config.yaml
Restart=always
RestartSec=10
StandardOutput=append:/var/log/orderbook-service.log
StandardError=append:/var/log/orderbook-service.log

[Install]
WantedBy=multi-user.target